# conf/tac_plus_init.conf
# tac_plus daemon
description "tac_plus daemon"
author      "Marc Huber"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
# Specify working directory
chdir /opt/tacacsgui
exec tac_plus.sh